# steamdeck-deployment/systemd/shader-predict-thermal.timer
[Unit]
Description=Thermal monitoring timer for Shader Prediction Compiler
Documentation=https://github.com/YourUsername/shader-prediction-compilation
Requires=shader-predict.service

[Timer]
# Check thermal state every 30 seconds
OnUnitActiveSec=30s
OnBootSec=60s

# Accuracy window for battery efficiency
AccuracySec=5s

# Persistent across reboots
Persistent=true

[Install]
WantedBy=timers.target